# quiz/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log

[program:api]
command=node /app/api/dist/index.js
directory=/app/api
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="4000"
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api_err.log

[program:ai]
command=/opt/venv/bin/python /app/ai/main.py
directory=/app/ai
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1",AI_SERVICE_PORT="8000"
stdout_logfile=/var/log/supervisor/ai.log
stderr_logfile=/var/log/supervisor/ai_err.log
